# com.clawfeed.digests.new.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawfeed.digests</string>
    <key>WorkingDirectory</key>
    <string>/Users/Zhuanz/Desktop/me/03-前端项目代码/xwzx-news/clawfeed</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/Zhuanz/.nvm/versions/node/v20.19.0/bin/node</string>
        <string>/Users/Zhuanz/Desktop/me/03-前端项目代码/xwzx-news/clawfeed/scripts/generate-digests.mjs</string>
        <string>--type=4h</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_PATH</key>
        <string>/Users/Zhuanz/.nvm/versions/node/v20.19.0/lib/node_modules</string>
    </dict>
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/Zhuanz/Desktop/me/03-前端项目代码/xwzx-news/clawfeed/logs/digests.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Zhuanz/Desktop/me/03-前端项目代码/xwzx-news/clawfeed/logs/digests.log</string>
</dict>
</plist>